# samples/conf/supervisor.conf
[group:grpc]
programs=samples

[program:samples]
directory=/data/php-grpc-server-protobuf/samples
command=php /data/php-grpc-server-protobuf/samples/server/server.php
autostart=true
autorestart=true
startretries=10
redirect_stderr=true
stdout_logfile=/tmp/supervisor.api.grpc.log
